# php-scaffold/image/phpfpm/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php-fpm
redirect_stderr=true

[include]
files=/etc/supervisor/conf.d/*.conf
